# services/gen-ai-hello-notes/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 * * * * /app/run-notes.sh >> /proc/1/fd/1 2>&1
